# deployment/tgparser.service
# Systemd service for Telegram Parser Web App
# Copy to: /etc/systemd/system/tgparser.service
# Enable: sudo systemctl enable tgparser
# Start: sudo systemctl start tgparser
# Status: sudo systemctl status tgparser
# Logs: sudo journalctl -u tgparser -f

[Unit]
Description=Telegram Comments Parser Web Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/Telescraper
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONUNBUFFERED=1"

# Production server with Gunicorn
ExecStart=/usr/bin/python3 -m gunicorn -w 2 -b 127.0.0.1:8000 --timeout 300 --access-logfile - --error-logfile - web_app:app

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tgparser

[Install]
WantedBy=multi-user.target
